# scripts/octi-pulpo.service
[Unit]
Description=Octi Pulpo — Swarm Coordination Brain
After=network.target redis.service
Wants=redis.service

[Service]
Type=forking
User=jared
WorkingDirectory=/home/jared/agentguard-workspace/octi-pulpo
EnvironmentFile=/home/jared/agentguard-workspace/octi-pulpo/.env
ExecStart=/bin/bash /home/jared/agentguard-workspace/octi-pulpo/scripts/start.sh
ExecStop=/bin/bash -c 'pkill -f "bin/octi-pulpo"; pkill -f "bin/octi-worker"; pkill -f "bin/octi-timer"'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
